# roles/openclaw_server/templates/openclaw.service.j2
[Unit]
Description=OpenClaw AI Assistant
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/openclaw serve --port {{ openclaw_server_port }}
Environment="OPENCLAW_LLM_PROVIDER={{ openclaw_server_llm_provider }}"
Environment="OPENCLAW_OLLAMA_URL={{ openclaw_server_ollama_url }}"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
