# scripts/mcp-anklume-dev.service
[Unit]
Description=AnKLuMe MCP Development Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/AnKLuMe
ExecStart=/usr/bin/python3 /root/AnKLuMe/scripts/mcp-anklume-dev.py --sse --port 9090
Restart=on-failure
RestartSec=5
Environment=ANKLUME_PROJECT_DIR=/root/AnKLuMe

[Install]
WantedBy=multi-user.target
