# host/boot/systemd/anklume-start.service
[Unit]
Description=anklume Setup — storage, framework & bootstrap
# Only run on live ISO boot (terminal-first: no display manager)
ConditionKernelCommandLine=boot=anklume
After=local-fs.target network-online.target incus.service
Wants=network-online.target
Requires=incus.service
Conflicts=getty@tty2.service

[Service]
Type=oneshot
ExecStart=/opt/anklume/scripts/start.sh
StandardInput=tty
StandardOutput=journal+console
StandardError=journal+console
TTYPath=/dev/tty2

[Install]
WantedBy=multi-user.target
